# debian/vyatta-cfg-vpn.postinst.in
#!/bin/bash

prefix=@prefix@
exec_prefix=@exec_prefix@
sysconfdir=@sysconfdir@
bindir=@bindir@
sbindir=@sbindir@

# remove init of daemons that we start/stop
for init in openswan ipsec setkey; do
  update-rc.d -f ${init} remove >/dev/null
done

# remove keys
rm -f /etc/ipsec.secrets
touch /etc/ipsec.secrets
chown root:root /etc/ipsec.secrets
chmod 600 /etc/ipsec.secrets
rm -f /etc/ipsec.d/private/localhost.localdomainKey.pem
rm -f /etc/ipsec.d/certs/localhost.localdomainCert.pem

# Local Variables:
# mode: shell-script
# sh-indentation: 4
# End:
